# cron/crontab
*/10 * * * * /app/cron/script.sh >> /var/log/cron.log 2>&1